# deployment/systemd/fte-briefing.service
[Unit]
Description=FTE CEO Briefing - Weekly report generation and email
After=network-online.target

[Service]
Type=oneshot
User=fte
Group=fte
WorkingDirectory=/opt/fte-employee
ExecStart=/opt/fte-employee/.venv/bin/fte briefing generate --email
StandardOutput=append:/var/log/fte/briefing.log
StandardError=append:/var/log/fte/briefing-error.log

# Environment variables (SMTP credentials, etc.)
EnvironmentFile=-/etc/fte/briefing.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/fte-employee /home/fte/AI_Employee_Vault /var/log/fte

[Install]
WantedBy=multi-user.target
